# code/extra-code/minecraft.service
[Unit]
Description=Minecraft Service

[Services]
Type=simple
User=minecraft
ExecStart=/usr/bin/minecraft


